# roles/jumbotron/files/fakta.timer
[Unit]
Description=Ring fakta bell

[Timer]
OnCalendar=*-*-* 20:45:00
Persistent=false

[Install]
WantedBy=timers.target
