# roles/jumbotron/files/fakta.service
[Unit]
Description=Ring fakta bell

[Service]
Type=oneshot
User=jumbotron
ExecStart=/home/jumbotron/fakta-bell.sh
