# roles/jumbotron/files/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin jumbotron --noclear %I $TERM
